# deployment/deb_HandyTimeServer/DEBIAN/postinst
chown root:root /etc/systemd/system/HandyTimeServer.service
chmod -x /etc/systemd/system/HandyTimeServer.service

systemctl --now enable HandyTimeServer
systemctl start HandyTimeServer
