# deployment/deb_LocalTimeServer/DEBIAN/postinst
chown root:root /etc/systemd/system/LocalTimeServer.service
chmod -x /etc/systemd/system/LocalTimeServer.service

systemctl --now enable LocalTimeServer
systemctl start LocalTimeServer
